# bas-scenarios/persistence-mechanisms.json
{"_meta":{"content_type":"bas_scenario","content_version":1,"content_hash":"sha256:a0589be6dda5960f30cbc6ddbdc88d238234d391352ba22264869c7f21737907","exported_at":1777321199,"exported_by":"Void SIEM","source":"api"},"id":"bas_persistence_mechanisms","name":"Persistence Mechanisms","description":"Simulates multiple persistence mechanisms including registry run key modification, new Windows service creation, scheduled task persistence, and WMI event subscription. Tests detection of techniques used to maintain access across reboots.","category":"persistence","mitre_techniques":["T1547.001","T1543.003","T1053.005","T1546.003"],"events_template":[{"event_type":"registry_change","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"Registry value set: HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\SecurityUpdate = C:\\ProgramData\\updater.exe"},{"event_type":"registry_change","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"Registry value set: HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\ConfigHelper = powershell.exe -enc SQBFAFgAKABOAGUAdwAtAE8AYgBqAGU="},{"event_type":"registry_change","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"Registry value set: HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Userinit = C:\\Windows\\system32\\userinit.exe,C:\\ProgramData\\backdoor.exe"},{"event_type":"service_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"SYSTEM","raw":"New service installed: Name=WindowsUpdateHelper DisplayName='Windows Update Helper Service' BinPath=C:\\ProgramData\\updater.exe StartType=AutoStart"},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"sc.exe create WindowsUpdateHelper binPath= C:\\ProgramData\\updater.exe start= auto DisplayName= \"Windows Update Helper Service\""},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"schtasks.exe /create /tn \"Microsoft\\Windows\\Maintenance\\ConfigCheck\" /tr C:\\ProgramData\\updater.exe /sc ONLOGON /ru SYSTEM /rl HIGHEST"},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"schtasks.exe /create /tn \"DailyHealthCheck\" /tr \"powershell.exe -ep bypass -f C:\\ProgramData\\health.ps1\" /sc DAILY /st 02:00"},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"wmic.exe /NAMESPACE:\\\\root\\subscription PATH __EventFilter CREATE Name=\"SecurityFilter\", EventNamespace=\"root\\cimv2\", QueryLanguage=\"WQL\", Query=\"SELECT * FROM __InstanceModificationEvent WITHIN 60 WHERE TargetInstance ISA 'Win32_PerfFormattedData_PerfOS_System'\""},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"wmic.exe /NAMESPACE:\\\\root\\subscription PATH CommandLineEventConsumer CREATE Name=\"SecurityConsumer\", CommandLineTemplate=\"C:\\ProgramData\\updater.exe\""},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"admin","raw":"wmic.exe /NAMESPACE:\\\\root\\subscription PATH __FilterToConsumerBinding CREATE Filter=\"__EventFilter.Name='SecurityFilter'\", Consumer=\"CommandLineEventConsumer.Name='SecurityConsumer'\""}],"expected_rules":["registry_persistence","new_service","scheduled_task_persistence","wmi_subscription"],"expected_severity":"high","timeout_seconds":120,"enabled":true}