# hunt-templates/scheduled-task-persistence.json
{"_meta":{"content_type":"hunt_template","content_version":1,"content_hash":"sha256:1944eaced0617e72df45f39c77bf0563d0a13c11e6cf004ede2ab495a2d37195","exported_at":1777318582,"exported_by":"Void SIEM","source":"api"},"id":"scheduled-task-persistence","name":"Scheduled Task Persistence","description":"Identify malicious scheduled tasks used for persistence","hypothesis_template":"Attackers are creating scheduled tasks to maintain persistence and execute payloads","queries":events from security_logs where raw contains '4698' AND event_type = 'authentication'
events from security_logs where raw contains 'schtasks' AND (raw contains '/create' OR raw contains '/run'),"mitre_tactics":Persistence,Execution,"mitre_techniques":T1053.005,"difficulty":"beginner","category":"persistence"}