# rules/cron-job-modified.json
{"_meta":{"content_type":"rule","content_version":1,"content_hash":"sha256:2b8304647159d2a1130f7ac87405bc7f868c610e5f4a470398f814175be0fd26","exported_at":1777321199,"exported_by":"Void SIEM","source":"api"},"id":"rule-1774740103184820525","name":"Cron Job Modified","description":"Detects modifications to cron jobs on Linux systems - potential persistence","query":"// Cron Job Modification\n// Changes to cron directories or crontab\nevents from endpoint\nwhere file_path contains \"/etc/cron\" or file_path contains \"/var/spool/cron\"\n  and event_type in (\"file_create\", \"file_modify\")\ngroup by hostname, username, file_path\nhaving count(*) >= 1\nwithin 300 seconds","language":"native","severity":"medium","enabled":true,"mitre_techniques":["T1053.003"],"created_at":1774740103,"updated_at":1774740103}