# rules/cron-job-modified.json
{"_meta":{"content_type":"rule","content_version":1,"content_hash":"sha256:28fa9ac89dfa2c38d4d3a4e15b357890689b93972d82c800860f479cd94d88d5","exported_at":1777801035,"exported_by":"Void SIEM","source":"api"},"id":"rule-1774740103184820525","name":"Cron Job Modified","description":"Detects modifications to cron jobs on Linux systems - potential persistence","query":"// Cron Job Modification\n// Changes to cron directories or crontab\nevents from endpoint\nwhere file_path contains \"/etc/cron\" or file_path contains \"/var/spool/cron\"\n  and event_type in (\"file_create\", \"file_modify\")\ngroup by hostname, username, file_path\nhaving count(*) >= 1\nwithin 300 seconds","language":"native","severity":"medium","enabled":true,"mitre_techniques":["T1053.003"],"created_at":1774740103,"updated_at":1777800907}